# workspaces/kubespray-workspace/supervisord-kubespray.conf
[program:octant]
directory=/home
command=/bin/sh -c " octant "
stderr_logfile = /var/log/octant/stderr.log
stdout_logfile = /var/log/octant/arastdout.log
logfile_maxbytes = 1024
