# workspaces/ansible-terraform-workspace/supervisord-infra.conf
[program:ara]
directory=/home
command=/bin/sh -c " ara-manage runserver 0.0.0.0:8029 "
stderr_logfile = /var/log/ara/stderr.log
stdout_logfile = /var/log/ara/arastdout.log
logfile_maxbytes = 1024
