# workspaces/ubuntu-workspace/supervisord.conf
[supervisord]
nodaemon=true
environment=LC_ALL='en_US.UTF-8',LANG='en_US.UTF-8'
childlogdir=/var/log/supervisord/
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=10MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)


[program:ttyd]
directory=/home
command=/bin/sh -c " ttyd -p 8026 /bin/zsh "

